# supervisor_config/mydjango.ini
[program:mydjango]
process_name=%(program_name)s
command = /root/anaconda3/envs/python3env/bin/python manage.py runserver 0.0.0.0:80
directory = /root/MeiTuan_spiders_DataAnalysis/mydjango
startsecs = 5
startretries = 3
user = root
redirect_stderr = true
stdout_logfile_maxbytes = 20MB
stdout_logfile_backups = 10
autostart=true
autorestart=true
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
stderr_logfile=/root/MeiTuan_spiders_DataAnalysis/mydjango/err.log
stdout_logfile=/root/MeiTuan_spiders_DataAnalysis/mydjango/out.log